# Astute/package_scripts/postinstall
#/bin/bash
launchDaemon="$3/Library/LaunchDaemons/com.astutegraphics.deployment.plist"
helperTool="com.astutegraphics.deployment"
helperToolPath="/Library/PrivilegedHelperTools/com.astutegraphics.deployment"
logPath="/tmp/com.astutegraphics.deployment.log"

# Copy the Privileged Helper Tool
/bin/cp -f "$3/Applications/Astute Manager.app/Contents/Library/LaunchServices/$helperTool" "$3"/"$helperToolPath"

# Check for and delete any pre-existing Launch Daemon
if [[ -e "$launchDaemon" ]]; then
	/bin/launchctl unload "$launchDaemon"
	rm -rf "$launchDaemon"
fi

# Create the new Launch Daemon
/usr/libexec/PlistBuddy -c "Add Label string" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set Label $helperTool" "$launchDaemon"

/usr/libexec/PlistBuddy -c "Add MachServices dict" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Add MachServices:$helperTool bool" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set MachServices:$helperTool true" "$launchDaemon"

/usr/libexec/PlistBuddy -c "Add Program string" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set Program $helperToolPath" "$launchDaemon"

/usr/libexec/PlistBuddy -c "Add ProgramArguments array" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Add ProgramArguments:0 string" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set ProgramArguments:0 $helperToolPath" "$launchDaemon"

/usr/libexec/PlistBuddy -c "Add StandardErrorPath string" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set StandardErrorPath $logPath" "$launchDaemon"

/usr/libexec/PlistBuddy -c "Add StandardOutPath string" "$launchDaemon"
/usr/libexec/PlistBuddy -c "Set StandardOutPath $logPath" "$launchDaemon"


/bin/launchctl load "$launchDaemon"
